# update-motd.d/20-quote
#!/bin/bash

quote=`fortune`

echo "$quote"
echo
